# crontab
0 * * * * bash /home/noiopen/wiki.noiopen.it/backup.sh
